# windows/evasion_registry_startup_shell_folder_modified.yml
# Rule version v1.0.1

- name: "Windows: Suspicious Startup Shell Folder Modification"
  severity: "High"
  description: "Identifies suspicious startup shell folder modifications to change the default Startup directory in order to bypass
                detections monitoring file creation in the Windows Startup folder."
  solution: "Isolate the involved host to prevent further post-compromise behavior.
      Determine the initial vector abused by the attacker and take action to prevent reinfection through the same vector.
      If the malicious file was delivered via phishing:
      Block the email sender from sending future emails. Block the malicious web pages.
      Remove emails from the sender from mailboxes. Consider improvements to the security awareness program.
      Run a full antimalware scan. This may reveal additional artifacts left in the system, persistence mechanisms, and malware components."
  category: "Persistence"
  tactic: "Registry Run Keys / Startup Folder"
  reference:
    - "https://attack.mitre.org/tactics/TA0003/"
    - "https://attack.mitre.org/techniques/T1547/001/"
  frequency: 60
  cache: 
    - allOf:
        - field: "logx.wineventlog.registry.data.strings"
          operator: "not regexp"
          value: "(C:\\\\ProgramData\\\\Microsoft\\\\Windows\\\\Start Menu\\\\Programs\\\\Startup|%ProgramData%\\\\Microsoft\\\\Windows\\\\Start Menu\\\\Programs\\\\Startup|%USERPROFILE%\\\\AppData\\\\Roaming\\\\Microsoft\\\\Windows\\\\Start Menu\\\\Programs\\\\Startup|C:\\\\Users\\\\(.+)\\\\AppData\\\\Roaming\\\\Microsoft\\\\Windows\\\\Start Menu\\\\Programs\\\\Startup)"
        - field: "logx.wineventlog.registry.path"
          operator: "regexp"
          value: "(HKLM\\\\Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Explorer\\\\User Shell Folders\\\\Common Startup|HKLM\\\\Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Explorer\\\\Shell Folders\\\\Common Startup|HKEY_USERS\\\\(.+)\\\\Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Explorer\\\\User Shell Folders\\\\Startup|HKEY_USERS\\\\(.+)\\\\Software\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Explorer\\\\Shell Folders\\\\Startup)"
      minCount: 1
      timeLapse: 60
      save:
        - field: "logx.wineventlog.event_data.SubjectUserName"
          alias: "SourceUser"
        - field: "logx.wineventlog.host.name"
          alias: "SourceHost"
        - field: "logx.wineventlog.event_data.SourceAddress"
          alias: "SourceIP"
        - field: "logx.wineventlog.event_data.SourcePort"
          alias: "SourcePort"
        - field: "logx.wineventlog.event_data.TargetUserName"
          alias: "DestinationUser"
        - field: "logx.wineventlog.event_data.DestAddress"
          alias: "DestinationIP"
        - field: "logx.wineventlog.event_data.DestPort"
          alias: "DestinationPort"
