# windows/microsoft_security_essentials_loading_definitions_failed.yml
# Rule version v1.0.1

- name: "Windows: Microsoft security essentials - Loading definitions failed. Using last good set"
  severity: "Medium"
  description: "Detects a failure in the process of loading Microsoft Security Essentials definitions. 
                Consequently, the system is using the latest valid definition available to maintain some level of protection, but it might be out of date compared to the most recent version."
  solution: "Make sure the computer has a stable internet connection so definition updates can be downloaded successfully.
            Try restarting the service to resolve temporary problems with loading definitions.
            If the problem persists, try performing a manual update by downloading the latest definitions from the official Microsoft Security Essentials website and installing them manually."
  category: "Execution"
  tactic: "Scheduled Task/Job"
  dataTypes: ["wineventlog"]
  reference:
    - "https://attack.mitre.org/tactics/TA0002/"
    - "https://attack.mitre.org/techniques/T1053/"
  frequency: 60
  cache:
    - allOf:
        - field: "logx.wineventlog.event_data.ProviderName"
          operator: "=="
          value: "Microsoft Antimalware"
        - field: "logx.wineventlog.event_id"
          operator: "=="
          value: 2004
      timeLapse: 60
      minCount: 1
      save:
        - field: "logx.wineventlog.event_data.SubjectUserName"
          alias: "SourceUser"
        - field: "logx.wineventlog.host.name"
          alias: "SourceHost"
        - field: "logx.wineventlog.host.ip.1"
          alias: "SourceIP"
        - field: "logx.wineventlog.event_data.IpPort"
          alias: "SourcePort"
        - field: "logx.wineventlog.event_data.TargetUserName"
          alias: "DestinationUser"
