# windows/microsoft_security_essentials_definitions_update_failed.yml
# Rule version v1.0.1

- name: "Windows: Microsoft security essentials - Definitions update failed"
  severity: "Medium"
  description: "Detects a problem with Microsoft Security Essentials definition update. 
                This alert indicates that the latest virus and malware definitions could not be successfully downloaded or applied."
  solution: "We recommend checking your internet connection: Make sure your computer has a stable internet connection so that updates can be downloaded correctly.
            Try to restart the service to resolve temporary issues related to updates.
            If the problem persists, try performing a manual update by downloading the latest definitions from the official Microsoft Security Essentials website and installing them manually."
  category: "Execution"
  tactic: "Scheduled Task/Job"
  dataTypes: ["wineventlog"]
  reference:
    - "https://attack.mitre.org/tactics/TA0002/"
    - "https://attack.mitre.org/techniques/T1053/"
  frequency: 60
  cache:
    - allOf:
        - field: "logx.wineventlog.event_data.ProviderName"
          operator: "=="
          value: "Microsoft Antimalware"
        - field: "logx.wineventlog.event_id"
          operator: "=="
          value: 2001
      timeLapse: 60
      minCount: 1
      save:
        - field: "logx.wineventlog.event_data.SubjectUserName"
          alias: "SourceUser"
        - field: "logx.wineventlog.host.name"
          alias: "SourceHost"
        - field: "logx.wineventlog.host.ip.1"
          alias: "SourceIP"
        - field: "logx.wineventlog.event_data.IpPort"
          alias: "SourcePort"
        - field: "logx.wineventlog.event_data.TargetUserName"
          alias: "DestinationUser"
